# launchd/com.local.claude.cost-alert.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.claude.cost-alert</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/leonardo.piedade/EDF/Personal/Github/claude-bridge/scripts/run_cost_alert.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/leonardo.piedade/EDF/Personal/Github/claude-bridge</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/leonardo.piedade</string>
    </dict>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/leonardo.piedade/.claude-bridge/cost-alert.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/leonardo.piedade/.claude-bridge/cost-alert.err</string>
</dict>
</plist>
