# launchd/com.local.claude-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.claude-bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/leonardo.piedade/EDF/Personal/Github/claude-bridge/run.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/leonardo.piedade/EDF/Personal/Github/claude-bridge</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/leonardo.piedade</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/Users/leonardo.piedade/.claude-bridge/launchd.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/leonardo.piedade/.claude-bridge/launchd.err</string>
</dict>
</plist>
